# shanming-backend/deploy/shanming-api.service
# 此文件为模板。请运行 deploy/install-systemd.sh 自动生成并安装到 systemd。
# install-systemd.sh 会根据实际部署路径（如 /opt/school/shanming-backend）写入正确配置。

[Unit]
Description=Shanming Mini-Program Backend (Node.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/school/shanming-backend
ExecStart=/usr/bin/node app.js
Restart=always

[Install]
WantedBy=multi-user.target
